# launchd/com.slack-mcp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.slack-mcp</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/mgandal/.cache/slack-mcp/start.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/mgandal/.cache/slack-mcp/launchd-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mgandal/.cache/slack-mcp/launchd-stderr.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/mgandal</string>
</dict>
</plist>
