# launchd/com.nanoclaw.paperpile-wiki.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nanoclaw.paperpile-wiki</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/mgandal/Agents/nanoclaw/scripts/paperpile-wiki/maintain.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/mgandal/Agents/nanoclaw/logs/paperpile-wiki.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mgandal/Agents/nanoclaw/logs/paperpile-wiki.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/mgandal</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Users/mgandal/.local/bin</string>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/mgandal/Agents/nanoclaw</string>
</dict>
</plist>
